# dayPeak/dayPeak.cron
0 12 * * * sudo /home/faisal/dayPeak/dayPeak.sh >> results
