# files/peak.cron
59 * * * * /usr/bin/hourPeak
59 12 * * * /usr/bin/dayPeak
