# files/peak.cron
59 * * * * /bin/hourPeak.sh
59 12 * * * /bin/dayPeak.sh
